# mac/com.alexassync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alexassync</string>

    <key>UserName</key>
    <string>YOUR_USERNAME</string>

    <key>ProgramArguments</key>
    <array>
        <string>/path/to/node</string>
        <string>/path/to/alexa-todoist-sync/shared/sync.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/path/to/alexa-todoist-sync</string>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/path/to/alexa-todoist-sync/logs/sync.log</string>

    <key>StandardErrorPath</key>
    <string>/path/to/alexa-todoist-sync/logs/sync-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
